# quorum/minecraft-witness.service
[Unit]
Description=Minecraft HA Quorum Witness Server
Documentation=https://github.com/Sommer2019/CommunicationNetwork
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/witness/witness-server.sh
Restart=always
RestartSec=3
StandardOutput=append:/var/log/minecraft-witness.log
StandardError=append:/var/log/minecraft-witness.log
# The witness must stay alive; give it time to bind the port on restart
TimeoutStartSec=10

[Install]
WantedBy=multi-user.target
