# scripts/minecraft-failover.service
[Unit]
Description=Minecraft HA Failover Watchdog
Documentation=https://github.com/Sommer2019/CommunicationNetwork
After=network-online.target glusterd.service redis.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/minecraft/scripts/failover-watchdog.sh
Restart=always
RestartSec=5
StandardOutput=append:/var/log/minecraft-failover.log
StandardError=append:/var/log/minecraft-failover.log

# Run as root so it can call systemctl and mount
User=root

[Install]
WantedBy=multi-user.target
